# systemd/nanoclaw.service
[Unit]
Description=NanoClaw - Personal Claude Assistant
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{PROJECT_ROOT}}
ExecStart={{NODE_PATH}} {{PROJECT_ROOT}}/dist/index.js
Restart=always
RestartSec=5

# Environment
Environment=PATH={{NODE_BIN_DIR}}:{{HOME}}/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME={{HOME}}
Environment=ASSISTANT_NAME=Nanomi

[Install]
WantedBy=multi-user.target
